# modules/site/templates/systemd.service.erb
[Unit]
Description=<%= @description %>
Requires=<%= @requires %>
After=<%= @after %>

[Service]
Type=<%= @type %>
User=<%= @user %>
ExecStart=<%= @execstart %>
