# modules/site/templates/systemd.timer.erb
[Unit]
Description=<%= @description %>

[Timer]
OnCalendar=<%= @oncalendar %>
Persistent=<%= @persistent %>

[Install]
WantedBy=timers.target
