# modules/site/files/pacaur.timer
[Unit]
Description=PacAur Autoupdater

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
