# modules/site/files/smart.timer
[Unit]
Description=SMART status mailer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
